# roles/ansible-role-certbot/templates/certbot_renewal.timer
[Unit]
Description=Try to renew certbot certs weekly

[Timer]
OnCalendar=weekly

[Install]
WantedBy=multi-user.target
