# roles/ansible-role-certbot/templates/certbot_renewal.service
[Unit]
Description=Renew certbot

[Service]
Type=simple
ExecStart=/bin/certbot-auto renew --nginx

[Install]
WantedBy=multi-user.target
