# roles/ansible-synapse-docker/templates/synapse_clear_history.timer
[Unit]
Description=Clear synapse messages every hour

[Timer]
OnCalendar={{ synapse_clear_history_oncalendar}}

[Install]
WantedBy=multi-user.target
